# cookbooks/mailcatcher/templates/default/mailcatcher.conf.erb
description "MailCatcher"

start on runlevel [2345]
stop on runlevel [!2345]

exec scl enable rh-ruby22 "/opt/rh/rh-ruby22/root/usr/local/bin/mailcatcher --ip=0.0.0.0 --smtp-port=25 --foreground" >> /var/log/mailcatcher.log 2>&1
respawn
